# services/microsocks.service
[Unit]
Description=microsocks SOCKS5 proxy for egress shielding
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/microsocks -i 127.0.0.1 -p 1080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
